# src/deb/control/postrm
#!/bin/sh

if [ "$1" = purge ]; then
	update-rc.d stool remove
fi
