# src/deb/control/postinst
#!/bin/sh
dir=/opt/ui/opt/tools/stool
test -d $dir
exists=$?
set -e

if [ "$1" = "configure" ]; then
    if id -u servlet >/dev/null 2>&1; then
        echo "servlet user already exists"
    else
        # TODO --no-create-home
        adduser --system --quiet --ingroup users --home /home/servlet servlet
    fi
fi


# TODO: there's no package dependency for this ...
export MAVEN_HOME=/opt/ui/opt/tools/wsdtools/maven/current
/opt/ui/bin/setup-stool -v -batch ${dir}
if test $exists -eq 1; then
  chown -R root:users ${dir}/downloads ${dir}/inbox ${dir}/logs ${dir}/service-wrapper ${dir}/conf ${dir}/tomcat ${dir}/wrappers
  sudo -u servlet ${dir}/bin/stool-raw.sh chown -stage overview
fi

if [ -x "/etc/init.d/stool" ]; then
  update-rc.d stool defaults >/dev/null
fi
