# src/deb/control/postinst
#!/bin/sh
set -e

user=stool
bin=/usr/share/stool
# TODO replace this check by some kind of configuration
if [ -d /opt/ui/opt/tools ] ; then
    home=/opt/ui/opt/tools/stool
    group=users
else
    home=/var/lib/stool
    group=stool
fi

if [ "$1" = "configure" ]; then
    echo  "configure home ${home}, group ${group}"
    if getent group ${group} >/dev/null 2>&1; then
        echo "re-using existing group: ${group}"
    else
        echo "creating group: ${group}"
        groupadd $group
        for file in /home/* ; do
            name=$(basename $file)
            if id -u ${name} >/dev/null 2>&1; then
                echo "add user to group: ${name}"
                usermod -a -G ${group} ${name}
            else
                echo "ignoring home directory: $file"
            fi
        done
    fi
    if id -u ${user} >/dev/null 2>&1; then
        echo "re-using existing user: ${user}"
    else
        echo "creating user: ${user}"
        if [ -d /home/${user} ] ; then
             echo "home directory for user ${user} already exists"
             exit 1
        fi
        adduser --system --ingroup ${group} --home /home/${user} ${user}
    fi
    if groups ${user} | cut --delimiter=: --fields=2 | grep -q -w ${group} ; then
        echo "user ${user} already in group ${group}"
    else
        echo "add user ${user} to group ${group}"
        usermod -a -G ${group} ${user}
    fi
    ln -s ${home} ${bin}/home

    # TODO: should be possible to generate overview without MAVEN_HOME
    export MAVEN_HOME=/not/used/for/overview/stage
    STOOL_BIN=${bin} sg ${group} "java -cp ${bin}/stool.jar net.oneandone.stool.setup.DebianHome ${home}"
    sudo -u ${user} ${bin}/stool-raw.sh chown -stage overview

    update-rc.d stool defaults
fi
