# src/deb/control/stool.cron.d
1 4 * * * servlet /opt/ui/opt/tools/stool-[stool.number]/bin/stool-raw.sh validate -all -email -stop > /dev/null
